# spec/mock/auditd.service
[Unit]
Description=Mock auditd service

[Service]
Type=simple
ExecStart=/bin/bash -c "while true ; do echo 'mock auditd running' ; /usr/bin/sleep 1 ; done"

[Install]
WantedBy=multi-user.target
